# defaults/scry_2.service
[Unit]
Description=Scry2 — MTGA companion
Documentation=https://github.com/ShawnMcCool/scry-2
After=network-online.target
Wants=network-online.target

[Service]
# `bin/scry_2 start` execs the BEAM in the foreground, so Type=exec is
# the right contract — systemd considers the unit started once the
# binary is exec()ed and watches it as the main process.
Type=exec
# Patched by scripts/install-linux to point at the staged release
# under $HOME/.local/lib/scry_2.
ExecStart=/path/to/scry_2/bin/scry_2 start
# `bin/scry_2 stop` sends a graceful shutdown via the BEAM RPC API.
ExecStop=/path/to/scry_2/bin/scry_2 stop
Restart=on-failure
RestartSec=5
# `mixed` sends SIGTERM to the main pid (so `bin/scry_2 stop` can run
# its graceful path) and SIGKILL to lingering children if they don't
# exit within TimeoutStopSec — matches what the dev unit uses.
KillMode=mixed
TimeoutStopSec=30
LimitNOFILE=65536

# The release wrapper consults RELEASE_TMP for the inter-process
# pipe; default it to a per-user runtime dir so we don't pollute /tmp.
Environment=RELEASE_TMP=%t/scry_2

[Install]
WantedBy=default.target
